# udev_rules/95-monitor-hotplug.rules
# Rule to trigger monitor-hotplug.sh on graphics change (plug/unplug)
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/local/bin/monitor-hotplug.sh"
